# .claude/skills/ai-news/config/com.ainews.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ainews.pipeline</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/julien.hovan/ClaudeWorkflows/AI-News-Reports/.claude/skills/ai-news/scripts/run_full_pipeline.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/julien.hovan/ClaudeWorkflows/AI-News-Reports</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/julien.hovan/ClaudeWorkflows/AI-News-Reports/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/julien.hovan/ClaudeWorkflows/AI-News-Reports/logs/launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/julien.hovan</string>
  </dict>
</dict>
</plist>
